# etc/systemd/system/wg-ddns-endpoint-monitor.timer.d/override.conf
# This file is a package-installed conffile.
# Edit the values below to change the detection cadence. The monitor derives its
# internal grace window from the timer unit, so no matching /etc/default change
# is required. Then run:
#   systemctl daemon-reload
#   systemctl restart wg-ddns-endpoint-monitor.timer

[Timer]
OnBootSec=
OnUnitActiveSec=
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s
RandomizedDelaySec=1s
